# bind/rootfs/etc/services.d/bind/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start bind
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

mkdir -p /data/cache
if [ ! -f /data/named.conf ]; then
    cp /etc/named.conf /data/named.conf
fi

## Run your program
exec /usr/sbin/named -g -c /data/named.conf
