# bind/rootfs/etc/services.d/bind/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start bind
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

## Run your program
exec /usr/sbin/named -f